# conf/crontab
# 每天 12:00（北京时间）执行一次每日领奖提醒任务
0 12 * * * curl -s https://sofun.online/cron/send-remind > /dev/null 2>&1

